# provisioning/roles/grimd/templates/grimd.service.j2
[Unit]
Description=Grimd

[Service]
Type=simple
ExecStart={{ grimd_install_location }}/grimd
User={{ builder }}

[Install]
WantedBy=multi-user.target